# scripts/etc/systemd/system/dash-auctions-match.timer
[Unit]
Description=Run dash-auctions-match.service every minutes

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
